# wayland/wayland.target
[Unit]
Description=Wayland
AllowIsolate=yes
