# picom/compositor.service
[Unit]
Description=picom - Graphical Compositor

[Service]
ExecStart=picom --experimental-backends

[Install]
WantedBy=xorg.target
