# picom/supervisord.conf
[program:compositor]
command=picom --experimental-backends
autostart=false
