# cal/gcal-sync.service
[Unit]
Description=Sync local calcurse appointments with remote Google Calendar events
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gcal-sync

[Install]
WantedBy=default.target
