# music/spotify/spotifyd.service
[Unit]
Description=spotify daemon

[Service]
Type=exec
ExecStart=%h/.local/bin/_spotifyd-wrapper
Restart=always
RestartSec=1s
