# utils/unclutter/unclutter.service
[Unit]
Description=unclutter - Idle Cursor Hiding
PartOf=xorg.target

[Service]
ExecStart=unclutter -idle 1

[Install]
WantedBy=xorg.target
