# deadd/notifications.service
[Unit]
Description=deadd - Dbus Notifications

[Service]
ExecStart=deadd-notification-center

[Install]
WantedBy=xorg.target
