# deadd/supervisord.conf
[program:notifications]
command=deadd-notification-center
autostart=false
autorestart=true
