# xorg/deadd/deadd.service
[Unit]
Description=deadd - Dbus Notifications
PartOf=xorg.target

[Service]
ExecStart=deadd-notification-center

[Install]
WantedBy=xorg.target
