# xorg/xorg.target
[Unit]
Description=X11
AllowIsolate=yes
