# dropbox/supervisord.conf
[program:dropbox]
command=%(ENV_HOME)s/.local/maestral/bin/maestral gui
autostart=true
