# screenlock/supervisord.conf
[program:screenlock]
command=xss-lock -n /usr/lib/xsecurelock/dimmer -l -- %(ENV_HOME)s/.local/bin/screenlock
autorestart=true
