# wayland/hyprland/hypridle/hypridle.service
[Unit]
Description=Wayland screenlock manager
PartOf=wayland.target

[Service]
ExecStart=hypridle

[Install]
WantedBy=wayland.target
