# luakit/adblock-update.service
[Unit]
Description=Update adblock filters

[Service]
Type=oneshot
ExecStart=%h/.local/bin/_adblock-update

[Install]
WantedBy=default.target
