# polybar/statusbar.service
[Unit]
Description=polybar - Status Bar

[Service]
ExecStart=polybar main
Restart=on-failure
