# polybar/supervisord.conf
[program:statusbar]
command=polybar main
autostart=false
