# utils/kanshi/kanshi.service
[Unit]
Description=Wayland adaptive displays
PartOf=Wayland.target

[Service]
ExecStart=kanshi

[Install]
WantedBy=wayland.target
