# unclutter/supervisord.conf
[program:unclutter]
command=unclutter -idle 1
autostart=false
